# thirdparty/dmredis/redis.service
[Unit]
Description=redis
Documentation=redis.service
After=network.target

[Service]
# cp redis.service /lib/systemd/system/redis.service
# cp redis-server /usr/local/bin/redis-server
# cp redis.conf /usr/local/bin/redis.conf

# rm -f /usr/lib/systemd/system/redis.service
# rm -f /usr/local/bin/redis-server
# rm -f /usr/local/bin/redis.conf

Type=simple
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/redis-server ./redis.conf
ExecStop=/bin/kill -QUIT $MAINPID
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
